# scripts/com.anima.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.anima.watchdog</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/michael/.anima/bin/watchdog</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/michael/Projects/iamclaudia-ai/anima</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/michael/.local/bin:/Users/michael/.bun/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>/Users/michael</string>
  </dict>

  <key>KeepAlive</key>
  <true/>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/michael/.anima/logs/watchdog-launchd.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/michael/.anima/logs/watchdog-launchd.log</string>
</dict>
</plist>
